# com.tradebot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradebot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>REPLACE_WITH_FULL_PATH/tradebot/start.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>3000</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>REPLACE_WITH_FULL_PATH/tradebot/logs/output.log</string>

    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_FULL_PATH/tradebot/logs/error.log</string>

    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_FULL_PATH/tradebot</string>
</dict>
</plist>
